# 32_file_io/2_day_homework/3/inittab
12323213213213
131223123122312
31231223123
131231231234124
123123123123123
adassaasdasdasdasdaa
